# contrib/cilium-docker.conf
start on starting cilium-net-daemon
stop on stopped cilium-net-daemon

exec cilium-docker
