# contrib/cilium-socket-proxy.conf
start on starting cilium-net-daemon
stop on stopped cilium-net-daemon

exec socat -d -d TCP-LISTEN:9000,fork UNIX:/var/run/cilium/cilium.sock
