# test/vtep/vxlan-responder.service
[Unit]
Description=Spark service

[Service]
ExecStart=/usr/local/bin/vxlan-responder.py --dport 8472 --outerip 172.18.0.1 --innerip 10.1.5.1 --vni 2 --sport 55555 --bridge br-22b28ede79c2

[Install]
WantedBy=multi-user.target
