# contrib/upstart/cilium-docker.conf
start on (started cilium)
stop on stopped cilium

exec cilium-docker --debug
