# contrib/upstart/cilium-etcd.conf
start on starting cilium

pre-start script
	if docker ps -f status=exited | grep cilium-etcd; then
		docker rm -f cilium-etcd
	fi
	if ! docker ps | grep cilium-etcd; then
        exec docker run -d \
            -v /usr/share/ca-certificates/:/etc/ssl/certs \
            -p 4001:4001 -p 2380:2380 -p 2379:2379 \
            --name cilium-etcd \
            quay.io/coreos/etcd:v3.1.0 \
            etcd -name etcd0 \
            -advertise-client-urls http://127.0.0.1:2379,http://127.0.0.1:4001 \
            -listen-client-urls http://0.0.0.0:2379,http://0.0.0.0:4001 \
            -initial-cluster-token etcd-cluster-1 \
            -initial-cluster-state new
	fi
end script

script
	sleepWhileAppIsUp(){
		while docker ps | grep "$1" >/dev/null; do
			sleep 2
		done
	}

	sleepWhileAppIsUp "cilium-etcd"
end script

post-stop script
	if docker ps -a | grep cilium-etcd; then
		docker rm -f cilium-etcd
	fi
end script

